# roles/plex/files/plex-media-server.service
[Unit]
Description=Dockerized Plex Media Server
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=docker start plex
ExecStop=docker stop plex
